# supervisor/conf.d/celery_flower.ini
[program:celery_flower]
command=celery -A LaatUI flower
directory=/Users/sunxinyang/Desktop/LaatUI/LaatUI
autostart=true
autorestart=unexpected
user=sunxinyang