# supervisor/conf.d/celery.ini
[program:celery]
command=celery -A LaatUI worker -l info
directory=/Users/sunxinyang/Desktop/LaatUI/LaatUI
autostart=true
autorestart=unexpected
user=sunxinyang
stdout_logfile=/Users/sunxinyang/Desktop/LaatUI/supervisor/var/log/celery-stdout.log
stderr_logfile=/Users/sunxinyang/Desktop/LaatUI/supervisor/var/log/celery-stderr.log