# supervisor/conf.d/laatui.ini
[program:LaatUI]
command=gunicorn LaatUI.wsgi -w 2 -k gthread -b 0.0.0.0:80
directory=/Users/sunxinyang/Desktop/LaatUI/LaatUI
autostart=true
autorestart=unexpected
user=sunxinyang
stdout_logfile=/Users/sunxinyang/Desktop/LaatUI/supervisor/var/log/myweb-stdout.log
stderr_logfile=/Users/sunxinyang/Desktop/LaatUI/supervisor/var/log/myweb-stderr.log
